# cache/supervisord.conf
[program:init1]
command=mkdir /mnt/data
user=xrootd

[program:init1.1]
command=chown xrootd:xrootd /mnt/data -R
user=root

[program:init2]
command=mkdir /xcache-root
user=root

[program:init2.1]
command=chown xrootd:xrootd /xcache-root
user=root

[program:init3]
command=mkdir /xcache-meta
user=root

[program:init3.1]
command=chown xrootd:xrootd /xcache-meta
user=root

[program:init4]
command=useradd cmsuser
user=root

[program:init5]
command=useradd cmsprod
user=root

[program:init6]
command=mkdir /run/xcache-auth/
user=root

[program:init6.1]
command=chown xrootd:xrootd /run/xcache-auth/
user=root

[program:init7]
command=sh usr/libexec/xcache/renew-proxy
user=xrootd

[program:xrootd-server]
command=xrootd -I v4 -c /etc/xrootd/xrootd.cfg -l /var/log/xrootd/xrootd.log
user=xrootd
autorestart=true
