# cache/supervisord.conf
[program:init]
command=mkdir /mnt/data
user=xrootd

[program:init2]
command=mkdir /xcache-root
user=root

[program:init3]
command=chown xrootd:xrootd /xcache-root
user=root

[program:init4]
command=mkdir /xcache-meta
user=root

[program:init5]
command=chown xrootd:xrootd /xcache-meta
user=root

[program:xrootd-server]
command=xrootd -I v4 -c /etc/xrootd/xrootd.cfg -l /var/log/xrootd/xrootd.log
user=xrootd
autorestart=true
