# cache/supervisord.conf
[program:init]
command=yum install -y voms-clients vo-client-lcmaps-voms
user=root

[program:init1]
command=mkdir /mnt/data
user=xrootd

[program:init1.1]
command=chown xrootd:xrootd /mnt/data -R
user=root

[program:init2]
command=mkdir /xcache-root
user=root

[program:init2.1]
command=chown xrootd:xrootd /xcache-root
user=root

[program:init3]
command=mkdir /xcache-meta
user=root

[program:init3.1]
command=chown xrootd:xrootd /xcache-meta
user=root

[program:init4]
command=useradd cmsuser
user=root

[program:init4]
command=mkdir /usr/libexec/xcache/
user=root

[program:xrootd-server]
command=xrootd -I v4 -c /etc/xrootd/xrootd.cfg -l /var/log/xrootd/xrootd.log
user=xrootd
autorestart=true
