# cluster/server/supervisord.conf
[program: server-cmsd]
command=cmsd -c /etc/xrootd/xrootd-third-party-copy.cfg -l /var/log/xrootd/cmsd.log
user=xrootd
autorestart=true

[program:xrootd-server]
command=xrootd -c /etc/xrootd/xrootd-third-party-copy.cfg -l /var/log/xrootd/xrootd.log
user=xrootd
autorestart=true

[program: telegraf]
command=telegraf --config /etc/telegraf/telegraf.conf
autorestart=true
